# cron.d/mrtg
* * * * *	root	/bin/bash /opt/mrtg/cron1m.sh 
*/5 * * * *	root	/bin/bash /opt/mrtg/cron5m.sh
